# scripts/autonomy/com.ebull.autonomy.supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent for the autonomy SUPERVISOR (walk-away-for-days mode).
  KeepAlive + RunAtLoad: starts on load and is restarted if it ever exits
  (crash, reboot, OOM). The supervisor itself loops sessions forever with
  usage-limit backoff, so launchd's only job is to keep ONE supervisor alive.
  Its internal lock prevents duplicates if launchd double-starts.

  Use this OR com.ebull.autonomy.plist (the simpler hourly run_loop) — not both.
  Install: see scripts/autonomy/setup.md
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ebull.autonomy.supervisor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/lukebradford/Dev/eBull/scripts/autonomy/supervisor.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/lukebradford/Dev/eBull</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <!-- Resurrect the supervisor if it exits for any reason; throttle restarts. -->
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>/tmp/ebull-autonomy-supervisor.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ebull-autonomy-supervisor.err.log</string>
</dict>
</plist>
